# System/etc/systemd/system/dgpu.service
[Unit]
Description=Disable Nvidia.

[Service]
ExecStart=/home/lordarcadius/LinuxStuff/Scripts/dgpu.sh

[Install]
WantedBy=multi-user.target
